# debian-serialkiller/files/supervisord/currenttime.conf
[program:currenttime]
command=python3 /usr/local/bin/currenttime.py
user=root
autostart=true
autorestart=true
