# debian-serialkiller/files/supervisord/nodejs.conf
[program:nodejs]
command=node /var/www/server.js
user=root
autostart=true
autorestart=true
