# debian-serialkiller/files/supervisord/weather_montpellier.conf
[program:weather_montpellier]
command=python3 /usr/local/bin/weather_montpellier.py
user=root
autostart=true
autorestart=true
